# debian/eta-check-optictouch.postinst
#!/bin/sh

set -e

systemctl enable eta-check-optictouch.service
systemctl start eta-check-optictouch.service

#DEBHELPER#
