# debian/eta-check-optictouch.prerm
#!/bin/sh

set -e

systemctl stop eta-check-optictouch.service
systemctl disable eta-check-optictouch.service
