# debian/eta-check-optictouch.prerm
#!/bin/sh

set -e

systemctl disable eta-check-optictouch.service
