# debian/eta-check-optictouch.postinst
#!/bin/sh

set -e

systemctl enable eta-check-optictouch.service
